# services/watchcat.timer
[Unit]
Description=Trigger Dorian Watchdog

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
